# root/etc/crontabs/root
#########################################################################
# 【砍掉的优化全合并】坚果云备份插件定时任务
#########################################################################
# 每日凌晨2点自动执行日常备份
0 2 * * * /usr/bin/jianguoyun-backup.sh auto_daily
# 每月1号凌晨3点自动执行全量备份
0 3 1 * * /usr/bin/jianguoyun-backup.sh auto_full
# 每天凌晨1点自动清理过期备份文件
0 1 * * * /usr/bin/jianguoyun-backup.sh auto_clean
